# service/abstrac-engine-stop.service
[Unit]
Description=Abstrac Rendering Engine

[Service]
ExecStart=systemctl stop abstrac-engine

[Install]
WantedBy=default.target
